# server/atii-sync.service
[Unit]
Description=ATII Sync Server - сервер синхронизации данных
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/atii/server
Environment="NODE_ENV=production"
Environment="PORT=3001"
Environment="DATA_DIR=/var/www/atii/server"
Environment="STATIC_DIR=/var/www/atii/build"
ExecStart=/usr/bin/node /var/www/atii/server/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
